# launchd/com.advai.crypto-bot.news-collector.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.advai.crypto-bot.news-collector</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/advai/crypto-bot/.venv/bin/python</string>
    <string>/Users/advai/crypto-bot/news_collector.py</string>
    <string>--db</string>
    <string>/Users/advai/crypto-bot/trend_bot.db</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/advai/crypto-bot</string>

  <key>RunAtLoad</key>
  <true/>

  <key>StartInterval</key>
  <integer>600</integer>

  <key>StandardOutPath</key>
  <string>/Users/advai/crypto-bot/news_collector.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/advai/crypto-bot/news_collector.err.log</string>
</dict>
</plist>
